# web/res/wiki-ci-dev.service.tmpl
[Unit]
Description=workspace-ci Wiki - Next.js Dev Server (HMR)
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=__PROJECT_ROOT__
Environment="PATH=__PROJECT_ROOT__/node_modules/.bin:__NODE_BIN__:/usr/local/bin:/usr/bin:/bin"
Environment="WIKI_DEV_PORT=__DEV_PORT__"
Environment="WIKI_DEV_HOST=0.0.0.0"
EnvironmentFile=%h/.config/wiki-ci-dev.env
ExecStart=__NODE__ scripts/dev-single.mjs
# Dev: only restart on failure, not on clean exit (Restart=always resurrected
# stale processes and fought manual stop/restart during development).
Restart=on-failure
RestartSec=5
MemoryHigh=4G
MemoryMax=6G
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wiki-ci-dev

[Install]
WantedBy=default.target
